# cor-T1053.005.xml
<group name="windows,sysmon,attack,persistence,execution,t1053.005,scheduled_task,">

  <!-- ============================= -->
  <!-- Baseline: Wazuh Default Rule  -->
  <!-- ============================= -->
  <!-- Rule ID 60228 (built-in): 
       "A scheduled task was created" 
       -> Windows Security Event ID 4698 -->
  <!-- We don’t redefine it; just build on top of it. -->

  <!-- ============================= -->
  <!-- Child Rules (fidelity)        -->
  <!-- ============================= -->

  <!-- Suspicious flag: /SC ONCE -->
  <rule id="200101" level="10">
    <if_sid>60228</if_sid>
    <regex field="win.eventdata.CommandLine">(?i)\/sc\s+once</regex>
    <description>Scheduled Task created with /SC ONCE (one-time execution) on $(win.system.computer). Command: $(win.eventdata.CommandLine)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>scheduled_task,suspicious_flag,</group>
  </rule>

  <!-- Suspicious flag: /TR launching interpreters -->
  <rule id="200102" level="11">
    <if_sid>60228</if_sid>
    <regex field="win.eventdata.CommandLine">(?i)\/tr\s+.*(cmd\.exe|powershell\.exe|wscript\.exe|mshta\.exe)</regex>
    <description>Scheduled Task created with /TR launching an interpreter on $(win.system.computer). Command: $(win.eventdata.CommandLine)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>scheduled_task,suspicious_target,</group>
  </rule>

  <!-- Suspicious TaskName: not under \Microsoft\Windows\ -->
  <rule id="200103" level="9">
    <if_sid>60228</if_sid>
    <negate>
      <regex field="win.eventdata.TaskName">(?i)\\Microsoft\\Windows\\</regex>
    </negate>
    <description>Scheduled Task created outside Microsoft\Windows path: $(win.eventdata.TaskName) on $(win.system.computer).</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>scheduled_task,non_standard_path,</group>
  </rule>

  <!-- ============================= -->
  <!-- Sysmon Process Creation Rule  -->
  <!-- ============================= -->

  <!-- schtasks.exe /create (Sysmon EventID 1) -->
  <rule id="200104" level="10">
    <if_group>process_creation</if_group>
    <regex field="win.eventdata.Image">(?i)\\schtasks\.exe$</regex>
    <regex field="win.eventdata.CommandLine">(?i)\/create</regex>
    <description>Sysmon: schtasks.exe used with /create on $(win.system.computer). Command: $(win.eventdata.CommandLine)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>scheduled_task,sysmon,</group>
  </rule>

  <!-- ============================= -->
  <!-- Correlation Rules             -->
  <!-- ============================= -->

  <!-- Correlation: Event Log 60228 + Sysmon schtasks.exe -->
  <rule id="200199" level="13">
    <if_sid>200104,60228</if_sid>
    <timeframe>300</timeframe> <!-- 5 min window -->
    <description>Confirmed Scheduled Task creation: seen in both Event Log (4698) and Sysmon (process creation). Command: $(win.eventdata.CommandLine)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>correlation,scheduled_task,high_confidence,</group>
  </rule>

  <!-- Correlation: Suspicious flags + Sysmon schtasks.exe -->
  <rule id="200200" level="14">
    <if_sid>200101,200102,200104</if_sid>
    <timeframe>300</timeframe>
    <description>High Confidence: Suspicious Scheduled Task (one-time or interpreter target) confirmed by Sysmon execution and Event Log.</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>correlation,scheduled_task,high_confidence,</group>
  </rule>

  
  <!-- Atomic 1: Suspicious file read inside PowerShell -->
  <rule id="210105" level="6">
    <if_sid>92027</if_sid> <!-- Sysmon PowerShell event -->
    <regex field="win.eventdata.CommandLine">(?i)System\.IO\.File\]\:\:ReadAllText</regex>
    <description>PowerShell suspicious file read (System.IO.File::ReadAllText API)</description>
    <group>powershell,file_access,hunting,</group>
  </rule>

  <!-- Atomic 2: PowerShell Invoke-CimMethod -->
  <rule id="210106" level="8">
    <if_sid>92027</if_sid>
    <regex field="win.eventdata.CommandLine">(?i)Invoke\-CimMethod</regex>
    <description>PowerShell executed Invoke-CimMethod (WMI invocation)</description>
    <group>powershell,wmi,hunting,</group>
  </rule>

  <!-- Atomic 3: PS_ScheduledTask class access -->
  <rule id="210107" level="9">
    <if_sid>92027</if_sid>
    <regex field="win.eventdata.CommandLine">(?i)PS_ScheduledTask</regex>
    <description>PowerShell accessed WMI class PS_ScheduledTask</description>
    <mitre><id>T1053.005</id></mitre>
    <group>powershell,wmi,scheduled_task,</group>
  </rule>

  <!-- Atomic 4: RegisterByXml method -->
  <rule id="210108" level="10">
    <if_sid>92027</if_sid>
    <regex field="win.eventdata.CommandLine">(?i)RegisterByXml</regex>
    <description>PowerShell invoked RegisterByXml method of PS_ScheduledTask</description>
    <mitre><id>T1053.005</id></mitre>
    <group>powershell,wmi,scheduled_task,</group>
  </rule>

  <!-- Correlation: Full malicious sequence -->
  <rule id="210201" level="14">
    <if_sid>210101,210102,210103,210104</if_sid>
    <timeframe>300</timeframe> <!-- 5 minutes -->
    <description>High-confidence: PowerShell abused WMI (PS_ScheduledTask.RegisterByXml) to create a scheduled task with XML payload</description>
    <mitre>
      <id>T1053.005</id>
      <tactic>Persistence</tactic>
      <tactic>Execution</tactic>
      <technique>Scheduled Task/Job: Scheduled Task</technique>
    </mitre>
    <group>attack,correlation,persistence,scheduled_task,powershell,wmi,high_confidence,</group>
  </rule>

   <!-- Rule 1: Registry hijack of .msc association -->
  <rule id="230101" level="10">
    <if_sid>60202</if_sid> <!-- Sysmon process creation, reg.exe -->
    <regex field="win.eventdata.Image">(?i)\\reg\.exe</regex>
    <regex field="win.eventdata.CommandLine">(?i)reg\s+add\s+.*HKCU\\Software\\Classes\\mscfile\\shell\\open\\command</regex>
    <description>Registry hijack attempt: .msc file association modified via reg.exe ($(win.eventdata.CommandLine))</description>
    <mitre>
      <id>T1546.015</id>
      <id>T1112</id>
    </mitre>
    <group>registry_event,hijack,persistence,</group>
  </rule>

  <!-- Rule 2: Suspicious scheduled task pointing to .msc (expanded set) -->
  <rule id="230102" level="9">
    <if_sid>60228</if_sid> <!-- Scheduled task creation -->
    <regex field="win.eventdata.CommandLine">(?i)(compmgmt\.msc|eventvwr\.msc|services\.msc|taskschd\.msc|fsmgmt\.msc)</regex>
    <regex field="win.eventdata.CommandLine">(?i)\/sc\s+onlogon</regex>
    <regex field="win.eventdata.CommandLine">(?i)\/rl\s+highest</regex>
    <description>Suspicious scheduled task created: hijacked .msc target at logon with elevated privileges ($(win.eventdata.CommandLine))</description>
    <mitre><id>T1053.005</id></mitre>
    <group>scheduled_task,persistence,privilege_escalation,msc_hijack,</group>
  </rule>

  <!-- Rule 3: Correlation - Full malicious sequence -->
  <rule id="230199" level="15">
    <if_sid>230101,230102</if_sid>
    <timeframe>600</timeframe> <!-- 10 minutes -->
    <description>High-confidence persistence: Registry hijack of .msc + scheduled task creation for UAC bypass (CompMgmt, EventVwr, Services, Taskschd, Fsmgmt)</description>
    <mitre>
      <id>T1546.015</id>
      <id>T1053.005</id>
      <tactic>Persistence</tactic>
      <tactic>Privilege Escalation</tactic>
      <tactic>Execution</tactic>
    </mitre>
    <group>correlation,persistence,scheduled_task,registry,high_confidence,msc_hijack,</group>
  </rule>
</group>

